# deploy/systemd/aegis.service
# Native host unit (Sprint 10 + Sprint 35 watchdog).
# Run: node src/host/main.ts (Node >= 24).
[Unit]
Description=Aegis personal agent host
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=notify
WatchdogSec=60
NotifyAccess=main
User=aegis
Group=aegis
WorkingDirectory=/opt/aegis
EnvironmentFile=-/etc/aegis/env
Environment=AEGIS_CONFIG=/opt/aegis/aegis.config.json
ExecStart=/usr/bin/node /opt/aegis/src/host/main.ts
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/aegis

[Install]
WantedBy=multi-user.target
